# root/etc/s6-overlay/s6-rc.d/svc-calibre-web-automated/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

export CALIBRE_DBPATH=/volume/calibre-web-automated/config

# If CWA_PORT_OVERRIDE is not set, but PORT is set, use PORT
if [ -z "$CWA_PORT_OVERRIDE" ] && [ -n "$PORT" ]; then
    export CWA_PORT_OVERRIDE="$PORT"
fi

#------------------------------------------------------------------------------------------------------------------------
#  Pre-flight database permission check
#------------------------------------------------------------------------------------------------------------------------

echo "[svc-calibre-web] Running pre-flight database permission check..."

# Check app.db is writable (as abc user)
if [[ -f "/volume/calibre-web-automated/config/app.db" ]]; then
    if s6-setuidgid abc test -w "/volume/calibre-web-automated/config/app.db"; then
        echo "[svc-calibre-web] app.db is writable"
    else
        echo "[svc-calibre-web] WARNING: app.db is NOT writable by abc user!"
        echo "[svc-calibre-web] Attempting to fix permissions..."
        chown abc:abc /volume/calibre-web-automated/config/app.db 2>/dev/null
        chmod 664 /volume/calibre-web-automated/config/app.db 2>/dev/null
    fi
fi

# Check metadata.db is writable (as abc user) - critical for book metadata updates
if [[ -f "/volume/calibre-web-automated/calibre-library/metadata.db" ]]; then
    if s6-setuidgid abc test -w "/volume/calibre-web-automated/calibre-library/metadata.db"; then
        echo "[svc-calibre-web] metadata.db is writable"
    else
        echo "[svc-calibre-web] WARNING: metadata.db is NOT writable by abc user!"
        echo "[svc-calibre-web] This will cause 'Database Error: attempt to write a readonly database' when updating book metadata."
        echo "[svc-calibre-web] Attempting to fix permissions..."
        chown abc:abc /volume/calibre-web-automated/calibre-library/metadata.db 2>/dev/null
        chmod 664 /volume/calibre-web-automated/calibre-library/metadata.db 2>/dev/null
        # Also fix the directory for SQLite journal files
        chown abc:abc /volume/calibre-web-automated/calibre-library 2>/dev/null
        chmod 775 /volume/calibre-web-automated/calibre-library 2>/dev/null
    fi
fi

# Fix permissions on all book directories (author folders and book folders)
# This is critical for editing book metadata which requires write access to book directories
echo "[svc-calibre-web] Fixing permissions on calibre-library contents..."
if [[ -d "/volume/calibre-web-automated/calibre-library" ]]; then
    # Try to fix ownership recursively (may fail on some volume mounts, that's OK)
    chown -R abc:abc /volume/calibre-web-automated/calibre-library 2>/dev/null || true
    # Set directory permissions to 775 (rwxrwxr-x) - allows abc user and group to write
    find /volume/calibre-web-automated/calibre-library -type d -exec chmod 775 {} \; 2>/dev/null || true
    # Set file permissions to 664 (rw-rw-r--) - allows abc user and group to write
    find /volume/calibre-web-automated/calibre-library -type f -exec chmod 664 {} \; 2>/dev/null || true
    echo "[svc-calibre-web] Calibre library permissions fixed"
fi

# Check config directory is writable (needed for SQLite journal/WAL files)
if [[ -d "/volume/calibre-web-automated/config" ]]; then
    if s6-setuidgid abc test -w "/volume/calibre-web-automated/config"; then
        echo "[svc-calibre-web] config directory is writable"
    else
        echo "[svc-calibre-web] WARNING: config directory is NOT writable!"
        chown abc:abc /volume/calibre-web-automated/config 2>/dev/null
        chmod 775 /volume/calibre-web-automated/config 2>/dev/null
    fi
fi

echo "[svc-calibre-web] Pre-flight check complete. Starting Calibre-Web..."

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${CWA_PORT_OVERRIDE:-8083}" \
        cd /app/calibre-web-automated s6-setuidgid abc python3 /app/calibre-web-automated/cps.py
